# supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3

[program:nextjs]
directory=/app/frontend
command=node_modules/.bin/next start -p 3000 -H 0.0.0.0
autostart=true
autorestart=true
stdout_logfile=/app/logs/nextjs.out.log
stderr_logfile=/app/logs/nextjs.err.log
environment=PORT=3000

[program:scheduler]
directory=/app
command=python scheduler.py
autostart=true
autorestart=true
stdout_logfile=/app/logs/scheduler.out.log
stderr_logfile=/app/logs/scheduler.err.log
